# ansible/roles/app/templates/hello.service.j2
# hello.service
# Generated by ansible at {{ now() }}

[Unit]
Description="Hello application instance on port {{ app_port }}"
Wants=network-online.target
After=network-online.target

[Service]
Type=simple
WorkingDirectory=/opt/{{ app_work_directory }}
ExecStart=/usr/bin/java -jar -Dserver.port={{ app_port }} /opt/{{ app_work_directory }}/app.jar
Restart=always

[Install]
WantedBy=default.target
